# daylight-calendar/rootfs/etc/services.d/daylight-calendar/finish
#!/usr/bin/execlineb -S0
# ==============================================================================
# Take down the S6 supervision tree when Daylight Calendar fails
# ==============================================================================
if { s6-test ${1} -ne 0 }
if { s6-test ${1} -ne 256 }

# Clean up any running processes
if [ -f "/run/chromium.pid" ]; then
  kill -15 $(cat /run/chromium.pid) || true
  rm -f /run/chromium.pid
fi

if [ -f "/run/openbox.pid" ]; then
  kill -15 $(cat /run/openbox.pid) || true
  rm -f /run/openbox.pid
fi

if [ -f "/run/xorg.pid" ]; then
  kill -15 $(cat /run/xorg.pid) || true
  rm -f /run/xorg.pid
fi

s6-svscanctl -t /var/run/s6/services 